# 99-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", RUN+="/home/afnarel/.monitor-detect/monitor-hotplug.sh"
